# supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=/usr/bin/Xvfb :1

[program:zeppelin]
priority=900
command=/zeppelin/bin/zeppelin.sh
